# supervisor.conf
[program:docx-to-pdf]
command=yarn start:production
directory=/root/pdf_generator
numprocs=1
autostart=true
autorestart=true