# supervisor.conf
[program:docx-to-pdf]
command=yarn start:production
directory=/docx-to-pdf
numprocs=1
autostart=true
autorestart=true